# munin/files/default/munin.cron
#
# cron-jobs for munin
#

MAILTO=root

* * * * *     munin test -x /usr/bin/munin-cron && /usr/bin/munin-cron
