# munin/files/default/munin.cron
#
# cron-jobs for munin, drop in /etc/cron.d on a Linux machine.
# MAILTO=root

MAILTO=root

*/5 * * * * munin if [ -x /usr/local/munin/bin/munin-cron ]; then /usr/local/munin/bin/munin-cron; fi
